# resources/yin_yang.service
[Unit]
Description=Automatic light and dark mode
After=suspend.target

[Service]
ExecStart=/usr/bin/yin-yang --systemd
