# resources/yin_yang.service
[Unit]
Description=Automatic light and dark mode

[Service]
ExecStart=/usr/bin/yin-yang --systemd
