# resources/yin_yang.timer
[Unit]
Description=Automatic light and dark mode

[Timer]
OnCalendar=08:00:00
OnCalendar=20:00:00
Persistent=true

[Install]
WantedBy=timers.target
